# templates/config/shared/saimprowler/saimprowler.timer
[Unit]
Description="Saimprowler Cronjob"

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Unit=saimprowler.service

[Install]
WantedBy=multi-user.target
